# root/etc/services.d/cadvisor/run
#!/usr/bin/with-contenv sh

echo "Starting cadvisor"
exec /opt/dataloop/cadvisor \
     --listen_ip="" --port=8080 >> /var/log/dataloop/cadvisor 2>&1

